# roles/secure-server/files/rc.local
#!/bin/sh
# https://bugzilla.redhat.com/show_bug.cgi?id=1490629
#
/usr/libexec/ufw/ufw-init start
